# etc/systemd/system/update-ablock-hosts.timer
[Unit]
Description=Timer for updating StevenBlack hosts daily.

[Timer]
OnBootSec=15s
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target